# gate@.service
[Unit]
Description=Gate 代理服务节点 (%i)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/gate-core run -c /etc/gate/%i.json
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
